# deploy/workoutdb-server.service
[Unit]
Description=WorkoutDB home server
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=exec
User=workoutdb
Group=workoutdb
WorkingDirectory=/opt/workoutdb/current
EnvironmentFile=/opt/workoutdb/shared/.env
ExecStart=/opt/workoutdb/current/.venv/bin/uvicorn workoutdb_server.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=3

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/opt/workoutdb/shared/db

[Install]
WantedBy=multi-user.target
